# deploy/systemd/timestd-metrology.target
[Unit]
Description=HF-TimeStd Metrology — All Channels
Documentation=https://github.com/mijahauan/hf-timestd
After=timestd-core-recorder.service
Wants=timestd-core-recorder.service

[Install]
WantedBy=multi-user.target
