# systemd/grape-upload-retry.service
[Unit]
Description=GRAPE Upload Retry — drain stuck PSWS uploads
Documentation=https://github.com/mijahauan/hf-timestd
# We expect this to run alongside / between grape-daily runs.  Fail
# fast and let the timer fire us again later if the queue is big.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=timestd
Group=timestd
WorkingDirectory=/var/lib/timestd
Environment=PATH=/opt/hf-timestd/venv/bin:/usr/bin:/bin

# `grape upload --resume` walks every <data-root>/upload/<YYYYMMDD>/
# directory still on disk (= cleanup didn't run = upload not yet
# confirmed), resets failed-status tasks back to pending, and drives
# the queue.  Exits 0 cleanly when there's nothing to do.
ExecStart=/opt/hf-timestd/venv/bin/python3 -m hf_timestd.cli grape upload --resume

StandardOutput=journal
StandardError=journal
SyslogIdentifier=grape-upload-retry

# Lighter than grape-daily — no decimation, just SFTP traffic.
CPUQuota=20%
MemoryMax=512M

[Install]
WantedBy=multi-user.target
