# systemd/grape-daily.timer
[Unit]
Description=GRAPE Daily Processing Timer
Documentation=https://github.com/HamSCI/hf-timestd
Requires=grape-daily.service

[Timer]
# Run daily at 01:00 UTC (after midnight, yesterday's data is complete)
OnCalendar=*-*-* 01:00:00 UTC

# Add randomized delay to avoid load spikes
RandomizedDelaySec=300

# Run on boot if missed (e.g., system was off)
Persistent=true

[Install]
WantedBy=timers.target
