# systemd/timestd-radiod-monitor.service
[Unit]
# Alert (journal breadcrumb via service-alert.sh; email if mailutils present)
# when this always-restart service finally gives up (StartLimit exhausted).
OnFailure=timestd-alert@%n.service
Description=HF Time Standard Radiod Health Monitor
Documentation=https://github.com/mijahauan/hf-timestd
After=network.target

[Service]
Type=simple
User=timestd
Group=timestd
WorkingDirectory=/opt/git/sigmond/hf-timestd

# Environment from file
# Operator identity (callsign, grid) from sigmond — CLIENT-CONTRACT §14.2
EnvironmentFile=-/etc/sigmond/coordination.env

EnvironmentFile=-/etc/hf-timestd/environment

# Run the radiod health monitor
ExecStart=/opt/git/sigmond/hf-timestd/venv/bin/python3 \
    /opt/git/sigmond/hf-timestd/scripts/monitor_radiod_health.py \
    /var/lib/timestd/state/radiod-status.json \
    10

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=timestd-radiod-monitor

# Process limits
LimitNOFILE=65536
Nice=5

[Install]
WantedBy=multi-user.target
