# systemd/timestd-vtec.service
[Unit]
Description=HF-TimeStd GNSS VTEC Monitor
Documentation=https://github.com/mijahauan/hf-timestd
After=network.target
Wants=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=notify
User=timestd
Group=timestd
# Operator identity (callsign, grid) from sigmond — CLIENT-CONTRACT §14.2
EnvironmentFile=-/etc/sigmond/coordination.env

EnvironmentFile=/etc/hf-timestd/environment
Environment=NETRC=/etc/hf-timestd/earthdata-netrc
Environment=MALLOC_ARENA_MAX=2
WorkingDirectory=/var/lib/timestd

# Watchdog: restart if no heartbeat for 60s (6 poll intervals at 10s each)
WatchdogSec=60
NotifyAccess=main

ExecStart=/opt/hf-timestd/venv/bin/python -u /opt/hf-timestd/scripts/live_vtec.py \
    --config /etc/hf-timestd/timestd-config.toml

# Resource limits
MemoryMax=400M
MemorySwapMax=0

# Always restart - VTEC data is important for ionospheric corrections
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=timestd-vtec

[Install]
WantedBy=multi-user.target
