# systemd/timestd-l2-calibration.service
[Unit]
# Alert (journal breadcrumb via service-alert.sh; email if mailutils present)
# when this always-restart service finally gives up (StartLimit exhausted).
OnFailure=timestd-alert@%N.service
Description=HF Time Standard L2 Calibration Service - Converts L1 to L2 Timing
Documentation=https://github.com/mijahauan/grape-recorder
After=timestd-metrology.service
Wants=timestd-metrology.service

[Service]
Type=notify
User=timestd
Group=timestd
# Operator identity (callsign, grid) from sigmond — CLIENT-CONTRACT §14.2
EnvironmentFile=-/etc/sigmond/coordination.env

EnvironmentFile=/etc/hf-timestd/environment
Environment=MALLOC_ARENA_MAX=2
WorkingDirectory=/var/lib/timestd

# Watchdog: restart if no heartbeat for 180s (3 poll intervals)
WatchdogSec=180

# Run L2 calibration service (reads station info and channels from config)
ExecStart=/opt/git/sigmond/hf-timestd/venv/bin/python -m hf_timestd.core.l2_calibration_service \
    --config /etc/hf-timestd/timestd-config.toml \
    --poll-interval 60.0 \
    --log-level INFO

# Resource limits - yield to radiod which handles real-time USB/FFT
Nice=-10
MemoryMax=600M
MemorySwapMax=0

Restart=always
RestartSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=timestd-l2-calibration

[Install]
WantedBy=multi-user.target
