# systemd/timestd-iono-reanalysis.service
[Unit]
Description=HF-TimeStd Ionospheric Reanalysis
Documentation=https://github.com/HamSCI/hf-timestd
After=network.target

[Service]
Type=oneshot
User=timestd
Group=timestd

# Run at lowest priority to avoid interfering with real-time services
Nice=19
IOSchedulingClass=idle

# Ensure output directories exist with correct permissions
# Operator identity (callsign, grid) from sigmond — CLIENT-CONTRACT §14.2
EnvironmentFile=-/etc/sigmond/coordination.env

ExecStartPre=+/bin/mkdir -p /var/lib/timestd/phase2/science/propagation_stats
ExecStartPre=+/bin/mkdir -p /var/lib/timestd/phase2/science/tec_reanalyzed
ExecStartPre=+/usr/bin/chown -R timestd:timestd /var/lib/timestd/phase2/science

ExecStart=/opt/git/sigmond/hf-timestd/venv/bin/python -m hf_timestd.core.ionospheric_reanalysis \
    --config /etc/hf-timestd/timestd-config.toml \
    --hours 1

# Timeout: should complete well within 5 minutes
TimeoutStartSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=timestd-iono-reanalysis

[Install]
WantedBy=multi-user.target
