# systemd/timestd-core-recorder.service
[Unit]
Description=HF Time Standard Core Recorder - Phase 1 RTP to Raw Buffer Archive
Documentation=https://github.com/mijahauan/hf-timestd
After=network-online.target ka9q-radio.service
Wants=network-online.target ka9q-radio.service

[Service]
Type=notify
User=timestd
Group=timestd

# Enable systemd watchdog
WatchdogSec=60
NotifyAccess=main

# Environment from file (set by install.sh)
EnvironmentFile=-/etc/hf-timestd/environment

# Working directory (production installation path)
WorkingDirectory=/opt/hf-timestd

# Core recorder V2 uses RadiodStream for RTP handling (no manual multicast)
ExecStart=/opt/hf-timestd/venv/bin/python3 -m hf_timestd.core.core_recorder_v2 \
    --config /etc/hf-timestd/timestd-config.toml
    
# Wait up to 5 minutes for startup (health check waits for 1st minute of data)
TimeoutStartSec=300

# Health check: Verify data is being written
ExecStartPost=/opt/hf-timestd/scripts/health-check-recorder.sh

# AGGRESSIVE restart policy - data loss is unacceptable
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

# Email alert on failure
OnFailure=timestd-alert@%n.service

# Logging to file and journal
StandardOutput=append:/var/log/hf-timestd/core-recorder.log
StandardError=inherit
SyslogIdentifier=timestd-core-recorder

# Resource limits - yield to radiod which handles real-time USB/FFT
Nice=5
MemoryMax=4G

# Security hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=read-only

# Paths we need write access to
ReadWritePaths=/var/lib/timestd /var/log/hf-timestd /dev/shm/timestd
ReadOnlyPaths=/etc/hf-timestd

# Network access for RTP multicast
PrivateNetwork=false

# Capabilities for multicast reception
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_RAW CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
