# systemd/timestd-radiod-monitor.service
[Unit]
Description=HF Time Standard Radiod Health Monitor
Documentation=https://github.com/mijahauan/hf-timestd
After=network.target radiod@ac0g-bee1-rx888.service
Wants=radiod@ac0g-bee1-rx888.service

[Service]
Type=simple
User=timestd
Group=timestd
WorkingDirectory=/opt/hf-timestd

# Environment from file
EnvironmentFile=-/etc/hf-timestd/environment

# Run the radiod health monitor
ExecStart=/opt/hf-timestd/venv/bin/python3 \
    /opt/hf-timestd/scripts/monitor_radiod_health.py \
    /var/lib/timestd/state/radiod-status.json \
    10

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=timestd-radiod-monitor

# Process limits
LimitNOFILE=65536
Nice=5

[Install]
WantedBy=multi-user.target
