# systemd/timestd-alert@.service
[Unit]
Description=HF-TimeStd Service Failure Alert for %i
After=%i.service

[Service]
Type=oneshot
# Operator identity (callsign, grid) from sigmond — CLIENT-CONTRACT §14.2
EnvironmentFile=-/etc/sigmond/coordination.env

ExecStart=/opt/git/sigmond/hf-timestd/scripts/service-alert.sh %i
