# systemd/timestd-hpps-watchdog.service
# /etc/systemd/system/timestd-hpps-watchdog.service
[Unit]
Description=HF-Timestd HPPS (T6 BPSK SHM) Watchdog
Documentation=https://github.com/mijahauan/hf-timestd
After=chrony.service
Wants=chrony.service

[Service]
Type=oneshot
User=root
Group=root

# Runs to completion in <2s on the happy path; only does work on
# the rare LastRx-exceeded path (which fires a systemctl restart).
#
# Threshold tuning (2026-05-21): each restart causes the BPSK
# calibrator to re-disambiguate against T4 (LAN GPS), and the
# resulting effective chain_delay can shift by hundreds of µs from
# the previous calibration.  That made an aggressive 120 s threshold
# actively harmful — restarts introduced larger discontinuities
# than the glitches they were trying to recover from.
#
# Current values:
#   HPPS_LASTRX_THRESHOLD_S=600  (10 min — only react to genuine outages)
#   HPPS_RESTART_COOLDOWN_S=1800 (30 min — prevent thrash even when fired)
#
# Today (2026-05-23) the chain_delay IS persisted across restarts AND
# disambig uses T5 (LB-1421 USB NMEA, GPS-truth-direct), so the
# original 120/300 rationale for tighter thresholds is moot.  Once a
# stable week of operation confirms the T5+persistence pipeline is
# resilient, this can revert to 120/300.
Environment=HPPS_LASTRX_THRESHOLD_S=600
Environment=HPPS_RESTART_COOLDOWN_S=1800
ExecStart=/opt/git/sigmond/hf-timestd/scripts/hpps-watchdog.sh

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hpps-watchdog
