# systemd/timestd-fusion.service
[Unit]
Description=HF-Timestd Multi-Broadcast Fusion (Chrony Feed)
Documentation=https://github.com/mijahauan/hf-timestd
After=timestd-l2-calibration.service
Requires=timestd-l2-calibration.service
StartLimitIntervalSec=300
StartLimitBurst=5
OnFailure=timestd-alert@%n.service

[Service]
Type=notify
User=timestd
Group=timestd

# Volatile state directory for the authority manager's contract file
# /run/hf-timestd/fusion_status.json — see METROLOGY.md §4.5.
# systemd creates /run/hf-timestd/ owned by timestd:timestd on start
# and removes it on stop.
RuntimeDirectory=hf-timestd
RuntimeDirectoryMode=0755

# Systemd Watchdog - service must send WATCHDOG=1 every 120 seconds
# The fusion service sends watchdog notifications in its main loop
# Increased from 30s to 120s to allow for HDF5 read operations on first cycle
WatchdogSec=120

# Note: Using Type=notify (changed from simple) to support watchdog
# The fusion service calls systemd_daemon.notify('READY=1') on startup
# and systemd_daemon.notify('WATCHDOG=1') in the main loop

# Operator identity (callsign, grid) from sigmond — CLIENT-CONTRACT §14.2
EnvironmentFile=-/etc/sigmond/coordination.env

Environment=MALLOC_ARENA_MAX=2

# Use jemalloc instead of glibc malloc — empirically 22% slower RSS
# growth on bee1 (2026-05-15 leak diagnostic, see commits ee90148 +
# 2f7a431 and project_hf_timestd_fusion_audit.md). The leak is
# primarily in C-extension internal state (h5py reader-per-cycle
# pattern), so jemalloc doesn't eliminate the growth — but its arena
# management is markedly better than glibc's for this workload.
# Falls back gracefully if libjemalloc isn't installed: ld.so logs
# one stderr warning and continues with the default allocator.
Environment=LD_PRELOAD=/usr/lib/x86_64-linux-gnu/libjemalloc.so.2

# Run as the installed module from the venv
ExecStart=/opt/hf-timestd/venv/bin/python -m hf_timestd.core.multi_broadcast_fusion \
    --data-root /var/lib/timestd \
    --interval 8.0 \
    --enable-chrony \
    --log-level INFO

# Ensure directories exist and have correct permissions
# Note: + prefix executes with full privileges (root) before dropping to User=timestd
ExecStartPre=+/bin/mkdir -p /var/lib/timestd/phase2/fusion
ExecStartPre=+/usr/bin/chown -R timestd:timestd /var/lib/timestd/phase2/fusion

# Stop chrony before removing SHM to prevent a race condition:
# During rapid fusion restarts, a previous ExecStartPost may restart chrony
# concurrently. Chrony calls shmget(IPC_CREAT|0600), recreating segments as
# root:0600 between our ipcrm and fusion's connect(). Stopping chrony first
# ensures the segments stay removed until fusion recreates them with 0666.
ExecStartPre=+/bin/sh -c 'systemctl stop chrony 2>/dev/null || true'

# Remove stale Chrony SHM segments with wrong permissions (root:0600)
# -M = remove by key; ignore errors if segments don't exist (|| true)
ExecStartPre=+/bin/sh -c 'ipcrm -M 0x4e545030 2>/dev/null || true'
ExecStartPre=+/bin/sh -c 'ipcrm -M 0x4e545031 2>/dev/null || true'

# Restart Chrony so it reattaches to the new SHM segments created by fusion.
# Fusion creates fresh SHM on startup; chrony holds stale segment references
# from the previous fusion instance and will show reach=0 until restarted.
# The + prefix runs as root (required for systemctl).
ExecStartPost=+/bin/sh -c 'sleep 3 && systemctl restart --no-block chrony 2>/dev/null || true'

# Resource limits - yield to radiod which handles real-time USB/FFT
Nice=-10
MemoryMax=800M
MemoryHigh=600M
MemorySwapMax=0

# Always restart - critical for time discipline
Restart=always
RestartSec=10

# Unified journal sink — see docs/DEBUGGING.md §2 for access patterns.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=timestd-fusion

[Install]
WantedBy=multi-user.target
