# systemd/timestd-web-ui.service
[Unit]
Description=TimeStd Recorder Web UI
Documentation=https://github.com/mijahauan/hf-timestd
After=network-online.target timestd-core-recorder.service
Wants=network-online.target

[Service]
Type=simple
# User/Group - timestd service user
User=timestd
Group=timestd
EnvironmentFile=-/etc/hf-timestd/environment

# Node.js production settings
Environment="NODE_ENV=production"
Environment="PORT=3000"

WorkingDirectory=/opt/hf-timestd/web-ui

ExecStart=/usr/bin/node monitoring-server-v3.js

Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
MemoryMax=512M

# Security hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=read-only

# Paths we need access to
ReadWritePaths=/var/lib/timestd /var/log/hf-timestd /dev/shm/timestd
ReadOnlyPaths=/etc/hf-timestd

# Network access
PrivateNetwork=false

# Logging
StandardOutput=append:/var/log/hf-timestd/web-ui.log
StandardError=inherit
SyslogIdentifier=timestd-web-ui

[Install]
WantedBy=multi-user.target
