# systemd/timestd-tsl3-watchdog.service
# /etc/systemd/system/timestd-tsl3-watchdog.service
[Unit]
Description=HF-Timestd TSL3 (T6 BPSK SHM) Watchdog
Documentation=https://github.com/mijahauan/hf-timestd
After=chrony.service
Wants=chrony.service

[Service]
Type=oneshot
User=root
Group=root

# Runs to completion in <2s on the happy path; only does work on
# the rare LastRx-exceeded path (which fires a systemctl restart).
#
# Threshold tuning (2026-05-21): each restart causes the BPSK
# calibrator to re-disambiguate against T4 (LAN GPS), and the
# resulting effective chain_delay can shift by hundreds of µs from
# the previous calibration.  That made an aggressive 120 s threshold
# actively harmful — restarts introduced larger discontinuities
# than the glitches they were trying to recover from.
#
# Current values:
#   TSL3_LASTRX_THRESHOLD_S=600  (10 min — only react to genuine outages)
#   TSL3_RESTART_COOLDOWN_S=1800 (30 min — prevent thrash even when fired)
#
# Revert to the prior 120/300 once chain_delay is persisted across
# restarts (see docs/HF-PPS-CHRONY-TUNING.md §5).
Environment=TSL3_LASTRX_THRESHOLD_S=600
Environment=TSL3_RESTART_COOLDOWN_S=1800
ExecStart=/opt/hf-timestd/scripts/tsl3-watchdog.sh

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tsl3-watchdog
