# systemd/timestd-pipeline-watchdog.service
[Unit]
Description=HF-TimeStd Pipeline Health Watchdog
Documentation=https://github.com/mijahauan/hf-timestd

[Service]
Type=oneshot
# Operator identity (callsign, grid) from sigmond — CLIENT-CONTRACT §14.2
EnvironmentFile=-/etc/sigmond/coordination.env

ExecStart=/opt/git/sigmond/hf-timestd/scripts/pipeline-watchdog.sh
User=root
Nice=19
