# systemd/timestd-web-api.service
[Unit]
Description=HF Time Standard Web API & UI
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=notify
User=timestd
Group=timestd
WorkingDirectory=/opt/hf-timestd/web-api
# NOTE: No PYTHONPATH set - the venv has hf_timestd installed as a proper package
# Setting PYTHONPATH to a non-existent or repo path causes import confusion
# Operator identity (callsign, grid) from sigmond — CLIENT-CONTRACT §14.2
EnvironmentFile=-/etc/sigmond/coordination.env

Environment="PATH=/opt/hf-timestd/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment=MALLOC_ARENA_MAX=2

# Watchdog: restart if no heartbeat for 60s (6 heartbeats at 10s each)
WatchdogSec=60
NotifyAccess=main

# Run with uvicorn directly
ExecStart=/opt/hf-timestd/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000

# Resource limits - yield to radiod which handles real-time USB/FFT
Nice=5
MemoryMax=600M
MemorySwapMax=0

# Always restart - web UI is important for monitoring
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=timestd-web-api

[Install]
WantedBy=multi-user.target
