# systemd/timestd-physics.service
[Unit]
Description=HF-Timestd Physics-Based Fusion (Science First)
Documentation=https://github.com/mijahauan/hf-timestd
After=timestd-analytics.service
Requires=timestd-analytics.service
PartOf=timestd-analytics.service

[Service]
Type=simple
User=timestd
Group=timestd

# Run as the installed module from the venv
ExecStart=/opt/hf-timestd/venv/bin/python -m hf_timestd.core.physics_fusion_service \
    --data-root /var/lib/timestd \
    --output /var/lib/timestd/phase2/fusion

# Ensure correct permissions
ExecStartPre=+/usr/bin/chown -R timestd:timestd /var/lib/timestd/phase2/fusion

# Restart on failure
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Email alert on failure
OnFailure=timestd-alert@%n.service

# Standard output logging
StandardOutput=append:/var/log/hf-timestd/physics.log
StandardError=append:/var/log/hf-timestd/physics.log

[Install]
WantedBy=multi-user.target
