# systemd/timestd-analytics.service
[Unit]
Description=HF Time Standard Analytics Service - Phase 2 Timing Analysis
Documentation=https://github.com/mijahauan/hf-timestd
After=network-online.target timestd-core-recorder.service
Wants=network-online.target
# Analytics depends on core recorder for data
Requires=timestd-core-recorder.service

[Service]
Type=forking
RemainAfterExit=yes

# User/Group - same as core recorder
User=hf-timestd
Group=hf-timestd

# Environment from file (set by install.sh)
EnvironmentFile=-/etc/hf-timestd/environment

# Working directory
WorkingDirectory=${TIMESTD_PROJECT}

# Analytics service launcher script
# This starts all 9 channel analyzers + multi-broadcast fusion as background processes
ExecStart=${TIMESTD_PROJECT}/scripts/timestd-analytics.sh start ${TIMESTD_CONFIG}
ExecStop=${TIMESTD_PROJECT}/scripts/timestd-analytics.sh stop

# Restart if the launcher fails
Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=timestd-analytics

# Resource limits (9 channels + fusion = ~10 processes)
# Yield to radiod which handles real-time USB/FFT
Nice=5
MemoryMax=4G

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

# Paths we need access to
ReadWritePaths=/var/lib/timestd /var/log/timestd /dev/shm/timestd
ReadOnlyPaths=/etc/hf-timestd

[Install]
WantedBy=multi-user.target
