# systemd/timestd-physics.service
[Unit]
Description=HF-Timestd Physics Service (L2 Interpretation)
Documentation=https://github.com/mijahauan/hf-timestd
After=timestd-l2-calibration.service
# Soft dependency only (was Requires=): physics is an asynchronous/batch
# pipeline and must never be torn down or gated by a metrology service.
# See docs/design/METROLOGY_PHYSICS_SPLIT.md.
Wants=timestd-l2-calibration.service
OnFailure=timestd-alert@%n.service

[Service]
# Type=simple (not notify): physics startup must neither gate, nor be gated
# by, anything in the metrology path. The service still pets the systemd
# watchdog via WATCHDOG= notifications (WatchdogSec below).
Type=simple
User=timestd
Group=timestd

# Watchdog: restart if no heartbeat for 120s (2 minutes)
WatchdogSec=120

# Operator identity (callsign, grid) from sigmond — CLIENT-CONTRACT §14.2
EnvironmentFile=-/etc/sigmond/coordination.env

Environment=MALLOC_ARENA_MAX=2

# Run as the installed module from the venv
ExecStart=/opt/git/sigmond/hf-timestd/venv/bin/python -m hf_timestd.core.physics_fusion_service \
    --data-root /var/lib/timestd \
    --output /var/lib/timestd/phase2/fusion

# Ensure physics-owned directories exist and have correct permissions.
# The chown is scoped to physics outputs (fusion, science) ONLY — chowning
# the whole phase2 tree would touch live L2 metrology channel files and can
# race the metrology SWMR writers. Physics must not write into, or modify,
# the metrology path (docs/design/METROLOGY_PHYSICS_SPLIT.md).
ExecStartPre=+/bin/mkdir -p /var/lib/timestd/phase2/fusion
ExecStartPre=+/bin/mkdir -p /var/lib/timestd/phase2/science
ExecStartPre=+/usr/bin/chown -R timestd:timestd /var/lib/timestd/phase2/fusion
ExecStartPre=+/usr/bin/chown -R timestd:timestd /var/lib/timestd/phase2/science

# Resource limits - yield to radiod which handles real-time USB/FFT
Nice=5
MemoryHigh=3G
MemoryMax=4G
MemorySwapMax=0

# Restart on failure
Restart=always
RestartSec=10

# Unified journal sink — see docs/DEBUGGING.md §2 for access patterns.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=timestd-physics

[Install]
WantedBy=multi-user.target
