# systemd/timestd-sqlite-parity.timer
[Unit]
Description=HF-Timestd SQLite parity check timer
Documentation=https://github.com/mijahauan/hf-timestd
Documentation=file:///opt/git/sigmond/hf-timestd/docs/HDF5-TO-SQLITE-MIGRATION.md

[Timer]
# Every 6 hours, calendar-anchored so restarts don't drift the schedule.
# OnBootSec=10min — first run shortly after boot to surface
# install-time misconfiguration quickly.
OnBootSec=10min
OnCalendar=*-*-* 00,06,12,18:00:00
Persistent=true
AccuracySec=60s

[Install]
WantedBy=timers.target
