# scripts/timestd-vtec.service
[Unit]
Description=HF-TimeStd GNSS VTEC Monitor
After=network.target

[Service]
Type=simple
User=mjh
WorkingDirectory=/home/mjh/git/hf-timestd
Environment=PYTHONPATH=src
ExecStart=/opt/hf-timestd/venv/bin/python3 -u scripts/live_vtec.py --config config/timestd-config.toml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
